# backend/infra/fstak-load-secrets.service
[Unit]
Description=Load fstak secrets from Secret Manager
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/bash /opt/fstak/backend/infra/load-secrets.sh
RemainAfterExit=true

[Install]
WantedBy=multi-user.target
